# dotfiles/etc/systemd/system/offlineimap.timer
[Unit]
Description=OfflineIMAP sync timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=offlineimap.service

[Install]
WantedBy=default.target
